# crontab
* * * * * cd ~/grive && grive
*/10 * * * * aptitude search '~U' > /tmp/listOfUpgrades
